# rootfs/mnt/etc/profile.d/qisda_path.sh
#!/bin/sh

#Hiko: This is for qisda's environment setting
#      and this must be put into /et/profile.d
#      and please make sure this file is exeutable
#      ex: -rwxr-xr-x root root


#Hiko {

#for WiFi
export LD_LIBRARY_PATH=/usr/local/sbin:/usr/local/lib

#Hiko }

